# cron.d/sensors
*/5 * * * * root /opt/jeedom-rpi-scripts/runSensors.sh
